# services/systemd/lab-camera.service
[Unit]
Description=NeverOT Lab Camera Monitor
After=network.target

[Service]
Type=simple
User=pi
ExecStart=/home/pi/lab-env/bin/python /home/pi/lab-pi-setup/services/camera_monitor.py
WorkingDirectory=/home/pi
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
